# rivalcfg-tray.service
[Unit]
Description=RivalCfg Tray App

[Service]
ExecStart=/usr/local/bin/rivalcfg-gui
Restart=on-failure

[Install]
WantedBy=default.target